# packaging/preremove.sh
#!/bin/bash

systemctl disable prometheus-slurm-exporter.service || /bin/true
